# my_wealth_tracker/rootfs/etc/services.d/frontend/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the frontend service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

bashio::log.info "${message:="Starting python frontend server..."}"

## Run your program
cd public || exit
python3 -m http.server 8080